# rust-mcp/debian/rust-mcp.service
[Unit]
Description=Odoo MCP Server - Model Context Protocol server for Odoo integration
Documentation=https://github.com/rachmataditiya/odoo-rust-mcp
After=network.target

[Service]
Type=simple

# Setup config directory before starting
ExecStartPre=/bin/mkdir -p /etc/rust-mcp
ExecStart=/usr/bin/rust-mcp-service --transport http --listen 127.0.0.1:8787
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

# Config and runtime directories
ConfigurationDirectory=rust-mcp
RuntimeDirectory=rust-mcp

# Working directory for static files
WorkingDirectory=/usr/share/rust-mcp

# Environment
Environment=RUST_LOG=info
Environment=HOME=/root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rust-mcp

[Install]
WantedBy=multi-user.target
